# files_to_include/etc/cron.d/batchcron
#!/bin/sh
01 01 * * * root /usr/bin/php -q -f /var/www/batch/cron.php >> /var/www/batchcron.log
