# files_to_include/etc/cron.d/elastic
#!/bin/sh
01 */1 * * * root /usr/bin/php -q -f /var/www/pages/tools/getcron.php >> /var/www/getcron.log
